# rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv bash

echo "Starting redis..."
redis-server "${REDIS_CONFIG_PATH}/redis.conf"
