# rootfs/etc/s6-overlay/s6-rc.d/redis-swarm/run
#!/command/with-contenv bash

/etc/redis-swarm/redis-swarm
