# rootfs/etc/s6-overlay/s6-rc.d/redis-sentinel-prepare/up
/docker-entrypoint.d/redis-sentinel-prepare.sh